# deploy/backup/backup-db.timer
# Fires the backup-db.service every day at 03:15 UTC.
#
# Pick a quiet hour for your traffic pattern. 03:15 is a good default
# for an EU-hosted service: deep night for Europe, end of US trading.

[Unit]
Description=Daily GEFO Postgres backup
Requires=backup-db.service

[Timer]
OnCalendar=*-*-* 03:15:00 UTC
# Catch up after a downtime — if the server was off at 03:15, run on next boot.
Persistent=true
# Slight randomisation so multiple servers don't hit S3 at the same instant.
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
